# rootfs/etc/crontabs/selfoss
*/5 * * * * /usr/sbin/logrotate /etc/logrotate.d/selfoss